# deploy/systemd/table-signup.service
[Unit]
Description=NJU Table Signup Flask Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/table
EnvironmentFile=/etc/table-signup.env
ExecStart=/opt/table/.venv/bin/gunicorn -c /opt/table/gunicorn.conf.py app:app
Restart=always
RestartSec=3
KillSignal=SIGQUIT
TimeoutStopSec=30
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
